# traffic-stats-collect.service
[Unit]
Description=Collect per-device sing-box traffic counters
After=sing-box.service traffic-stats-nft.service
Wants=traffic-stats-nft.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/traffic-stats/collector.py
User=root
Group=trafficstats
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/traffic-stats
