# traffic-stats-collect.timer
[Unit]
Description=Run the sing-box traffic collector every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
RandomizedDelaySec=5s
AccuracySec=1s
Unit=traffic-stats-collect.service

[Install]
WantedBy=timers.target
